# servers/ai/remote/watchdog.service
[Unit]
Description=Watchdog WebSocket Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=watchdog
ExecStart=/opt/watchdog/.venv/bin/python /opt/watchdog/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
