# servers/ai/deploy/ping.service
[Unit]
Description=MQTT Ping Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ping-svc
ExecStart=/opt/ping/.venv/bin/python /opt/ping/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
